# modules/redis/files/redis-server.conf
#
# Redis
#

description "Redis nosql DB"

start on runlevel [2345]
stop on runlevel [!2345]

exec redis-server

