# modules/kibana/files/kibana.conf
#
# Kibana
#

description "Kibana log spelunker"

start on runlevel [2345]
stop on runlevel [!2345]

setuid kibana

exec kibana

